# roles/docnow_role/templates/tweet-loader.service.j2
[Unit]
Description=DocNow Tweet Loader

[Service]
User={{ docnow_user }}
WorkingDirectory={{ docnow_app_root }}
ExecStart=/usr/bin/npm run tweet-loader
Restart=always

[Install]
WantedBy=multi-user.target
