# roles/docnow_role/templates/url-fetcher.service.j2
[Unit]
Description=DocNow Stream Loader

[Service]
User={{ docnow_user }}
WorkingDirectory={{ docnow_app_root }}
ExecStart=npm run url-fetcher

[Install]
WantedBy=multi-user.target